# scraper-legacy/crontab
# run at 10:30pm everyday
48 10 * * * cd / && /usr/local/bin/python pisspricer-scraper scrape-all > /var/logs/logfile 2>&1
