# pricescraper/crontab
# Container time is 13 hours behind NZD time
0 3 * * * cd / && /usr/local/bin/python /start.py crawl newworld >> /var/logs/logfile 2>&1
